# deploy/systemd/tracet-caddy.service
[Unit]
Description=TraceT (Caddy)
After=tracet.service
WantedBy=tracet.service
PartOf=tracet.service

[Service]
Type=simple
WorkingDirectory=/home/debian/TraceT/deploy
ExecStart=/usr/bin/podman compose --profile prod up --force-recreate --remove-orphans --exit-code-from caddy caddy
ExecStop=/usr/bin/podman compose --profile prod down caddy
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
